# roles/calibre/templates/calibre-import.path.j2
[Unit]
Description=Calibre File Import

[Path]
DirectoryNotEmpty={{ calibre_import_dir }}

[Install]
WantedBy=multi-user.target
