# roles/sway/templates/swayidle.timer.j2
[Unit]
Description=Swayidle after user logs in

[Timer]
OnStartupSec=20sec
Persistent=false

[Install]
WantedBy=timers.target

